# srt-maker.service
[Unit]
Description=SRT-Maker Transcription Service
After=network.target

[Service]
Type=simple
User=luka
WorkingDirectory=/mnt/hdd/Development/SRT-Maker
ExecStart=/mnt/hdd/Development/SRT-Maker/venv/bin/uvicorn main:app --host 127.0.0.1 --port 8765
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
